# salt/qubes-air/remotevm/files/autossh-tunnel.service
# autossh-qubesair@.service —— 部署到本地 Relay (sys-relay-*) 的 systemd 用户/系统单元
# =====================================================================
# 零入站方案落地: 用 autossh 维持一条【出站】持久 SSH 到 Remote-Relay。
#   - 正向调用 (本地->远端): transport 脚本经此连接 (ControlMaster 复用) 转发。
#   - 反向调用 (远端->本地): 同一连接的 `ssh -R` 反向端口转发 (config.template 里的
#     RemoteForward) 把远端回环端口转回本 Relay 回环 sshd(2222)。
#
# 家庭 NAT 后的远端主机因此【不需要开任何公网入站端口】—— 连接方向永远是 Relay 主动出站。
#
# 实例名 = 远端原始名 remote_name (与 QubesDB /remote/<local> 一致), 例:
#     systemctl start autossh-qubesair@dev
#
# 持久化注意 (Qubes AppVM): 本单元与其依赖 (autossh 包、~/.ssh/*) 必须经 bind-dirs 或
# 模板持久化, 否则 Relay 重启即丢。见 salt/qubes-air/remotevm/relay.sls。
# =====================================================================
[Unit]
Description=Qubes Air autossh 出站持久隧道 -> Remote-Relay (%i)
After=qubes-network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# %i = 远端原始名 (ssh Host)
Environment=AUTOSSH_GATETIME=0
Environment=AUTOSSH_PORT=0
# -M 0 关闭 autossh 自带监控端口, 改用 ssh 的 ServerAlive* (见 config.template) 判活。
# -N 不执行远端命令 (只维持隧道 + 端口转发)。-F 指定受控 config (含 RemoteForward)。
ExecStart=/usr/bin/autossh -M 0 -N -F %h/.ssh/config %i
Restart=always
RestartSec=10
# 以 Relay 内普通用户身份运行 (不需要 root)
User=user

[Install]
WantedBy=multi-user.target

# =====================================================================
# 待真机确认:
#   [B4] autossh 包在目标模板 (fedora-*) 里的包名与路径 (/usr/bin/autossh)。
#   [B5] Qubes AppVM 里用 system 单元还是 user 单元更稳 (本模板按 system 单元 + User=user)。
#   [B6] %h 在 systemd 单元里是否展开为 User 的家目录 (通常是 root 家目录, 若用 User=user
#        需改成绝对路径 /home/user/.ssh/config)。见 relay.sls 里的实际渲染。
# =====================================================================
